# deploy/systemd/portfolio-backup-prune.service
[Unit]
Description=Portfolio verified backup 7/4/6 reachability pruning
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
EnvironmentFile=/etc/portfolio/backup.env
UnsetEnvironment=BACKUP_AGE_IDENTITY AGE_SECRET_KEY AGE_SECRET_KEY_FILE COS_SECRET_ID COS_SECRET_KEY SMTP_PASSWORD
ExecStart=/usr/bin/bash /usr/local/libexec/portfolio/backup-dispatch.sh prune
TimeoutStartSec=2h
UMask=0077
Nice=15
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectClock=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectHostname=true
ProtectProc=invisible
ProcSubset=pid
RestrictSUIDSGID=true
RestrictRealtime=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
CapabilityBoundingSet=
AmbientCapabilities=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
ReadWritePaths=/var/backups/portfolio
ReadOnlyPaths=/opt/portfolio
ReadOnlyPaths=/usr/local/libexec/portfolio/backup-dispatch.sh
ReadOnlyPaths=/etc/portfolio/backup.env
ReadOnlyPaths=/etc/portfolio/rclone-backup-pruner.conf
InaccessiblePaths=/var/run/docker.sock
InaccessiblePaths=/var/lib/docker/volumes/portfolio-local-media/_data
InaccessiblePaths=/etc/portfolio/postgres.env
InaccessiblePaths=/etc/portfolio/portfolio.env
InaccessiblePaths=/etc/portfolio/rclone-media-reader.conf
InaccessiblePaths=/etc/portfolio/rclone-backup-uploader.conf
InaccessiblePaths=/etc/portfolio/rclone-backup-verifier.conf
InaccessiblePaths=/etc/portfolio/msmtp-backup.conf
