# deploy/systemd/portfolio-backup.timer
[Unit]
Description=Nightly Portfolio encrypted backup schedule

[Timer]
# 18:30 UTC is 02:30 Asia/Hong_Kong on the following calendar day.
OnCalendar=*-*-* 18:30:00 UTC
RandomizedDelaySec=15m
Persistent=true
AccuracySec=1m
Unit=portfolio-backup.service

[Install]
WantedBy=timers.target
